# systemd/overlay-audit.service
[Unit]
Description=Audit and correct overlay mounts
After=local-fs.target postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/bin/php /srv/app/www/sync/scripts/audit_overlays.php
User=root
WorkingDirectory=/srv/app/www/sync
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=overlay-audit

[Install]
WantedBy=multi-user.target
